# daemon/auto-patcher-cron
#! /usr/bin/env sh

. /etc/rc.subr

name="auto-patcher"
start_cmd="${name}_start"
stop_cmd=":"

patcher_start()
{
	echo "Nothing started."
}

load_rc_config $name
run_rc_command "$1"
